# systemd/mergerfs-pool.service
[Unit]
Description=drive-nmonit mergerfs Local Drive Pool
After=local-fs.target
Wants=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
# The actual ExecStart is generated by setup-mergerfs.sh to include detected drives
ExecStart=/usr/bin/mergerfs -o defaults,allow_other,category.create=mfs,minfreespace=1G,fsname=drive-nmonit-pool /mnt/drives/*/ /mnt/local-pool
ExecStop=/bin/umount -l /mnt/local-pool
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
